# s6/etc/s6/service.d/urxvtd/run
#!/bin/execlineb -P
fdmove -c 2 1
sudo --user=leon
urxvtd --quiet --opendisplay

