# s6/etc/s6/service.d/agetty-tty6/run
#!/bin/execlineb -P
fdmove -c 2 1
/sbin/agetty --autologin leon 38400 tty6 linux

