# s6/etc/s6/service.d/crond/run
#!/bin/execlineb -P
redirfd -w 1 /dev/null
fdmove 2 1
busybox crond -f
