# s6/etc/s6/service.d/mocp/run
#!/bin/execlineb -P
sudo --user=leon
mocp --server --foreground
