# s6/etc/s6/service.d/agetty-tty1/run
#!/bin/execlineb -P
fdmove -c 2 1
/sbin/agetty --noclear 38400 tty1 linux
